# scripts/abp_sync.cron
# Auto Banana Product - OneDrive Sync
# Runs every 7 days at 3:00 AM
0 3 */7 * * /opt/1panel/docker/compose/auto_banana_product/scripts/sync_to_onedrive.sh >> /var/log/abp_sync.log 2>&1
